# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=httpd-foreground
