# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=apache2ctl -DFOREGROUND
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
killasgroup=true
stopasgroup=true

[program:clamd]
command=clamd -c /etc/clamav/clamd.conf

[program:freshclam]
command=freshclam -d --config-file=/etc/clamav/freshclam.conf
