# cron/jessie-unoconv/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f
autostart=true
autorestart=true

[program:unoconv]
command=unoconv --listener
autostart=true
autorestart=true
user=www-data
